# rootfs/etc/services.d/nordvpn/finish
#!/bin/bash

nordvpn disconnect
nordvpn logout
rm -rf /run/nordvpn