# rootfs/etc/services.d/nordvpn/run
#!/bin/bash

if [[ ! -d /run/nordvpn ]]; then
  mkdir -m 0770 /run/nordvpn
fi

s6-notifyoncheck /usr/sbin/nordvpnd